# bot/systemd/vpnbanana-accounting.timer
[Unit]
Description=Run vpnbanana accounting every 5 min

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
