# scripts/reality-watcher.service
[Unit]
Description=REALITY SNI watcher (pick fastest SNI from pool every 30 min)
After=network-online.target x-ui.service
Wants=network-online.target

[Service]
Type=simple
# Полный замер пула раз в 30 мин (секунды); изменить: Environment=WATCH_INTERVAL_SEC=3600
Environment=WATCH_INTERVAL_SEC=1800
ExecStart=/usr/local/bin/reality-failover.sh watch
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
