# scripts/reality-watcher.service
[Unit]
Description=REALITY SNI watcher (pick fastest SNI from pool every 30 min)
After=network-online.target x-ui.service
Wants=network-online.target

[Service]
Type=simple
# Полный замер пула раз в 30 мин (секунды); изменить: Environment=WATCH_INTERVAL_SEC=3600
# Пул: приоритет sni-rotation-pool.txt (см. build-sni-rotation-pool.sh), иначе sni-candidates.txt
# Свой файл: Environment=CANDIDATES_FILE=/path/to/pool.txt
# После смены SNI: Profile-Update-Interval (часы) и Announce в БД 3x-ui
# Environment=SUB_UPDATES_HOURS=1
# Environment=BUMP_SUB_ANNOUNCE=1
# Environment=BUMP_SUB_ON_KEEP=1
# Environment=RESTART_XUI_ON_KEEP=0
# Environment=PROBE_PARALLEL=40
# Push (ntfy или Telegram — см. REALITY-FAILOVER.md)
# Telegram по умолчанию из панели 3x-ui (settings: tgBotToken, tgBotChatId); ниже — override
# Environment=NOTIFY_URL=https://ntfy.sh/secret-topic
# Environment=NOTIFY_TITLE=VPN
# Environment=TELEGRAM_BOT_TOKEN=
# Environment=TELEGRAM_ADMIN_CHAT_ID=
# Environment=TELEGRAM_CLIENT_CHAT_ID=
# Environment=TELEGRAM_CHAT_ID=
# Environment=SWITCH_MIN_IMPROVE_MS=50
# Environment=NOTIFY_ON_KEEP=0
Environment=WATCH_INTERVAL_SEC=1800
ExecStart=/usr/local/bin/reality-failover.sh watch
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
